# .config/systemd/user/usb-backup.path
[Unit]
Description=USB Backup Trigger

[Path]
PathExists=/media/%u/SYNTH_126
Unit=usb-backup.service

[Install]
WantedBy=default.target
